# config/crontab
*/1 * * * * test -f /var/www/.cron && php /var/www/.cron
